# 09_04_WebStorageArchivePattern/crontab
*/1  *  *  *  * root logrotate --state /var/run/log_rotate.state /etc/logrotate.d/myapp >>/var/log/logrotate 2>&1
